# playbooks/file-system/mount-name-set.mount
[Unit]
Description=file-storage

[Mount]
What=
Where=
Type=ext4
Options=defaults,noexec

[Install]
WantedBy=multi-user.target
